# etc/supervisord.d/grafana.ini
[program:grafana]
directory=/opt/grafana/
command=/opt/qnib/grafana/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
